# templates/keycloak.service.erb
[Unit]
Description=Keycloak service
After=network-online.target
Wants=network-online.target

[Service]
Type=idle
SyslogIdentifier=keycloak
<% if scope['keycloak::service_environment_file'] -%>
EnvironmentFile=<%= scope['keycloak::service_environment_file'] %>
<% end -%>
EnvironmentFile=<%= scope['keycloak::admin_env'] %>
<% if scope['keycloak::java_opts'] -%>
  <%- if scope['keycloak::java_opts'].is_a?(Array) -%>
Environment='JAVA_OPTS_APPEND=<%= scope['keycloak::java_opts'].join(' ') %>'
  <%- else -%>
Environment='JAVA_OPTS_APPEND=<%= scope['keycloak::java_opts'] %>'
  <%- end -%>
<% end -%>
Environment='JAVA_HOME=<%= scope['keycloak::java_home'] %>'
User=<%= scope['keycloak::user'] %>
Group=<%= scope['keycloak::group'] %>
ExecStart=<%= scope['keycloak::service_start_cmd'] %>
# TODO: remove once upgraded from Keycloak 25 to 26
<% unless (scope['keycloak::features'] || []).include?('persistent-user-sessions') -%>
ExecStartPost=-/usr/bin/rm -f <%= scope['keycloak::login_config'] %>
<% end -%>
TimeoutStartSec=600
TimeoutStopSec=600
SuccessExitStatus=0 143

[Install]
WantedBy=multi-user.target
